# configuration_app/systemd/khaudioconfig.service
[Unit]
Description=scheduler for SIP calls and mp3 streams
Requires=multi-user.target
After=multi-user.target

[Service]
Type=simple
ExecStart=/var/lib/configuration_app/app.py

[Install]
WantedBy=multi-user.target